# nodo-remoto/systemd/camera-capture.service
[Unit]
Description=Camera Capture — FFmpeg MPEG-TS segmenter
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi_damh/camera
ExecStart=/usr/bin/python3 nodo-remoto/capture/capturer.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
