# nodo-remoto/systemd/camera-sync.service
[Unit]
Description=Camera Sync Scheduler — nightly archive sender
After=network.target wg-quick@wg0.service
Wants=wg-quick@wg0.service

[Service]
Type=simple
User=pi_damh
WorkingDirectory=/home/pi_damh/camera
ExecStart=/usr/bin/python3 nodo-remoto/sync/sync_scheduler.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
